# atlantis-code-server.service
[Unit]
Description=code-server (VS Code in browser — port 5001)
After=network.target

[Service]
Type=simple
Environment=PATH=/home/viktor/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/sh -c 'exec code-server --bind-addr 127.0.0.1:5001 --auth none'
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
