# atlantis-server.service
[Unit]
Description=Atlantis OS — HTTP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/viktor/library/projects/main
ExecStart=/usr/bin/python3 /home/viktor/library/projects/main/server.py
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
